# configuration/supervisord/riak-cs.conf
[program:riak-cs]
command=riak-cs attach
priority=3