# src/configuration/supervisord/riak.conf
[program:riak]
command=riak attach
priority=1