# src/configuration/supervisord/stanchion.conf
[program:stanchion]
command=stanchion attach
priority=2